# installer/templates/supervisor.conf.j2
[program:{{ app_name }}]
user = {{ app_user }}
stdout_logfile = {{ app_virtualenv_dir }}/logs/gunicorn_supervisor.log
redirect_stderr = true
command = {{ app_virtualenv_dir }}/bin/gunicorn_run
