# installer/templates/supervisor.conf.j2
[program:{{ project_app_name }}]
user = gunicorn
stdout_logfile = /venv/logs/gunicorn_supervisor.log
redirect_stderr = true
command = /venv/bin/run_gunicorn.sh
